# php-fpm-apache2/s6-files/etc/s6-overlay/s6-rc.d/apache2/run
#!/command/execlineb -P
with-contenv

apachectl -D FOREGROUND